# services/trinity/systemd/omega-summary.timer
[Unit]
Description=OMEGA Summary Timer

[Timer]
OnBootSec=30s
OnUnitActiveSec=5m
Unit=omega-summary.service

[Install]
WantedBy=timers.target
